# packages/app-misc/actkbd/files/systemd/actkbd.service
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=Dynamic DNS updater client

[Service]
Type=forking
User=root
ExecStart=
